# bash/profile.d/00-personal.sh
#
# Personal environment variables (always enabled via bash/setup)
#

export SCRIPTS_DIR="$HOME/.local/bin/scripts"
export GITHUB_USER="$USER"
export REPOS_DIR="$HOME/repos"
