# bash/profile.d/10-local-bin.sh
#
# User-local binaries and scripts (always enabled via bash/setup)
#

: ${SCRIPTS_DIR:="$HOME/.local/bin/scripts"}

path_prefix "$SCRIPTS_DIR" "$HOME/.local/bin"
