# bash/profile.d/android.sh
#
# Android SDK (enabled via android/install)
#

export ANDROID_HOME="/opt/android-sdk"
export NDK_HOME="$ANDROID_HOME/ndk/27.1.12297006"

path_prefix "$ANDROID_HOME/platform-tools"
