# bash/profile.d/bun.sh
#
# Bun runtime (enabled via bun/install)
#

export BUN_INSTALL="$HOME/.bun"

path_prefix "$BUN_INSTALL/bin"
