# bash/profile.d/opencode.sh
#
# OpenCode CLI (enabled via opencode/install)
#

path_prefix "$HOME/.opencode/bin"
