# bash/profile.d/wget.sh
#
# Wget configuration (enabled via wget/setup)
#

: ${XDG_CONFIG_HOME:=$HOME/.config}

export WGETRC="$XDG_CONFIG_HOME/wget/wgetrc"
